# docs/automating_letter.service.example
# Example systemd service configuration for AutomatingLetter
# Place this file at: /etc/systemd/system/automating_letter.service

[Unit]
Description=Automating Letter Flask App
After=network.target

[Service]
Type=exec
User=shobbak
Group=shobbak
WorkingDirectory=/home/shobbak/AutomatingLetter
Environment=PATH=/home/shobbak/AutomatingLetter/venv/bin
ExecStart=/home/shobbak/AutomatingLetter/venv/bin/gunicorn -w 4 -b 0.0.0.0:5000 --timeout 120 --access-logfile /home/shobbak/AutomatingLetter/logs/access.log --error-logfile /home/shobbak/AutomatingLetter/logs/error.log app:app
Restart=always
RestartSec=3

# Environment variables
Environment=OPENAI_API_KEY=your_openai_api_key_here
Environment=GOOGLE_DRIVE_FOLDER_ID=your_drive_folder_id_here
Environment=FLASK_ENV=production
Environment=DEBUG_MODE=false
Environment=LOG_LEVEL=INFO

[Install]
WantedBy=multi-user.target

# Commands to set up and start the service:
# sudo systemctl daemon-reload
# sudo systemctl enable automating_letter.service  
# sudo systemctl start automating_letter.service
# sudo systemctl status automating_letter.service
